# testing/btest/Baseline/tools.systemd-generator.test-all-options/..dir1.zeek-proxy@.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Proxy %i
SourcePath=<...>/config1
After=zeek-setup.service
After=zeek-logger@.service
StopPropagatedFrom=zeek.target
StartLimitIntervalSec=0

[Service]
SyslogIdentifier=zeek-proxy-%i
Type=exec
User=zeek
Group=zeek
WorkingDirectory=<...>/proxy-%i
ReadWritePaths=<...>/proxy-%i
CPUAffinity=7,8
MemoryMax=1024M
Environment=PATH=<...>/bin:<...>/bin:<...>/bin:/bin
Environment=ZEEKPATH=<...>/generated-scripts:<...>/zeek:<...>/policy:<...>/site:<...>/builtin-plugins
Environment=CLUSTER_NODE=proxy-%i
Environment=GLOBAL_ENV=1
Environment=LD_PRELOAD=<...>/libjemalloc.so
Environment=I_AM_PROXY=1
Environment=MORE_PROXY=2
ExecStart=<...>/zeek policy<...>/broker Broker::disable_ssl=T policy<...>/systemd-generator loaded-by-proxy-1 loaded-by-proxy-2 local tuning/json-logs
Slice=zeek-proxies.slice
Restart=always
RestartSec=1
