# testing/btest/Baseline/tools.systemd-generator.test-all-options/..dir1.zeek-manager.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Manager
SourcePath=<...>/config1
After=zeek-setup.service
After=zeek-logger@.service
StopPropagatedFrom=zeek.target
StartLimitIntervalSec=0

[Service]
SyslogIdentifier=zeek-manager
Type=exec
User=zeek
Group=zeek
WorkingDirectory=<...>/manager
ReadWritePaths=<...>/manager
CPUAffinity=1,2,3,4,5,6,7,8,9,10,11,12
MemoryMax=1024M
Environment=PATH=<...>/bin:<...>/bin:<...>/bin:/bin
Environment=ZEEKPATH=<...>/generated-scripts:<...>/zeek:<...>/policy:<...>/site:<...>/builtin-plugins
Environment=CLUSTER_NODE=manager
Environment=GLOBAL_ENV=1
Environment=LD_PRELOAD=<...>/libjemalloc.so
Environment=I_AM_MANAGER=1
Environment=MORE_MANAGER=2
ExecStart=<...>/zeek policy<...>/broker Broker::disable_ssl=T policy<...>/systemd-generator loaded-by-manager-first loaded-by-manager-second local tuning/json-logs
Slice=zeek-manager.slice
Restart=always
RestartSec=1
