# testing/btest/Baseline/tools.systemd-generator.test-all-options/..dir1.zeek-logger@.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Logger %i
SourcePath=<...>/config1
After=zeek-setup.service
StopPropagatedFrom=zeek.target
StartLimitIntervalSec=0

[Service]
SyslogIdentifier=zeek-logger-%i
Type=exec
User=zeek
Group=zeek
WorkingDirectory=<...>/logger-%i
ReadWritePaths=<...>/logger-%i <...>/var
CPUAffinity=0,1
MemoryMax=512M
Environment=PATH=<...>/bin:<...>/bin:<...>/bin:/bin
Environment=ZEEKPATH=<...>/generated-scripts:<...>/zeek:<...>/policy:<...>/site:<...>/builtin-plugins
Environment=CLUSTER_NODE=logger-%i
Environment=GLOBAL_ENV=1
Environment=LD_PRELOAD=<...>/libjemalloc.so
Environment=I_AM_LOGGER=1
Environment=MORE_LOGGER=2
ExecStart=<...>/zeek policy<...>/broker Broker::disable_ssl=T policy<...>/systemd-generator loaded-by-logger-1 loaded-by-logger-2 local tuning/json-logs
Slice=zeek-loggers.slice
Restart=always
RestartSec=1
