# testing/btest/Baseline/tools.systemd-generator.test-zeek-systemd-generator/normal-dir.zeek-setup.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Setup
SourcePath=<...>/config1
PartOf=zeek.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
Nice=0
MemoryMax=
User=
Group=
WorkingDirectory=
ExecStart=mkdir -p <...>/generated-scripts
ExecStart=<...>/zeek-cluster-layout-generator -L 3 -P 2 -W 4 -p 20000 -a 127.0.2.1 -m 30000 -b 10.0.0.1 -o <...>/cluster-layout.zeek
ExecStart=mkdir -p <...>/zeek
ExecStart=chown test-zeek-user:test-zeek-group <...>/zeek
ExecStart=mkdir -p <...>/log-queue
ExecStart=chown test-zeek-user:test-zeek-group <...>/log-queue
ExecStart=mkdir -p <...>/manager
ExecStart=chown test-zeek-user:test-zeek-group <...>/manager
ExecStart=mkdir -p <...>/logger-1
ExecStart=chown test-zeek-user:test-zeek-group <...>/logger-1
ExecStart=mkdir -p <...>/logger-2
ExecStart=chown test-zeek-user:test-zeek-group <...>/logger-2
ExecStart=mkdir -p <...>/logger-3
ExecStart=chown test-zeek-user:test-zeek-group <...>/logger-3
ExecStart=mkdir -p <...>/proxy-1
ExecStart=chown test-zeek-user:test-zeek-group <...>/proxy-1
ExecStart=mkdir -p <...>/proxy-2
ExecStart=chown test-zeek-user:test-zeek-group <...>/proxy-2
ExecStart=mkdir -p <...>/worker-1
ExecStart=chown test-zeek-user:test-zeek-group <...>/worker-1
ExecStart=mkdir -p <...>/worker-2
ExecStart=chown test-zeek-user:test-zeek-group <...>/worker-2
ExecStart=mkdir -p <...>/worker-3
ExecStart=chown test-zeek-user:test-zeek-group <...>/worker-3
ExecStart=mkdir -p <...>/worker-4
ExecStart=chown test-zeek-user:test-zeek-group <...>/worker-4
RemainAfterExit=yes
