# testing/btest/Baseline/tools.systemd-generator.test-zeek-systemd-generator/normal-dir.zeek-worker@.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Worker %i
SourcePath=<...>/config1
After=zeek-setup.service
After=zeek-manager.service
After=zeek-logger@.service
After=zeek-proxy@.service
Requires=zeek-setup.service
PartOf=zeek.target
StartLimitIntervalSec=0

[Service]
Type=exec
Nice=-19
MemoryMax=1G
User=test-zeek-user
Group=test-zeek-group
WorkingDirectory=<...>/worker-%i
ReadWritePaths=<...>/worker-%i 
CapabilityBoundingSet=CAP_NET_RAW
AmbientCapabilities=CAP_NET_RAW
Environment=PATH=<...>/bin:<...>/bin:<...>/bin:/bin
Environment=ZEEKPATH=<...>/generated-scripts:<...>/zeek:<...>/policy:<...>/site:<...>/builtin-plugins
Environment=CLUSTER_NODE=worker-%i
ExecStart=<...>/zeek -i ${INTERFACE} policy<...>/systemd-generator ignore_checksums=T the-cluster-backend -B plugin-Zeek-Cluster_Backend_ZeroMQ
Slice=zeek-workers.slice
Restart=always
RestartSec=1
