# testing/btest/Baseline/tools.systemd-generator.test-zeek-systemd-generator/normal-dir.zeek-logger@.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Logger %i
SourcePath=<...>/config1
After=zeek-setup.service
Requires=zeek-setup.service
PartOf=zeek.target
StartLimitIntervalSec=0

[Service]
SyslogIdentifier=zeek-logger-%i
Type=exec
Nice=-3
MemoryMax=2001000K
User=test-zeek-user
Group=test-zeek-group
WorkingDirectory=<...>/logger-%i
ReadWritePaths=<...>/logger-%i <...>/var 
Environment=PATH=<...>/bin:<...>/bin:<...>/bin:/bin
Environment=ZEEKPATH=<...>/generated-scripts:<...>/zeek:<...>/policy:<...>/site:<...>/builtin-plugins
Environment=CLUSTER_NODE=logger-%i
ExecStart=<...>/zeek policy<...>/systemd-generator ignore_checksums=T the-cluster-backend -B plugin-Zeek-Cluster_Backend_ZeroMQ
Slice=zeek-loggers.slice
Restart=always
RestartSec=1
