# testing/btest/Baseline/tools.systemd-generator.test-zeek-systemd-generator/normal-dir.zeek-proxy@.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Proxy %i
SourcePath=<...>/config1
After=zeek-setup.service
After=zeek-logger@.service
Requires=zeek-setup.service
PartOf=zeek.target
StartLimitIntervalSec=0

[Service]
Type=exec
Nice=-2
MemoryMax=2000000000
User=test-zeek-user
Group=test-zeek-group
WorkingDirectory=<...>/proxy-%i
ReadWritePaths=<...>/proxy-%i 
Environment=PATH=<...>/bin:<...>/bin:<...>/bin:/bin
Environment=ZEEKPATH=<...>/generated-scripts:<...>/zeek:<...>/policy:<...>/site:<...>/builtin-plugins
Environment=CLUSTER_NODE=proxy-%i
ExecStart=<...>/zeek policy<...>/systemd-generator ignore_checksums=T the-cluster-backend -B plugin-Zeek-Cluster_Backend_ZeroMQ
Slice=zeek-proxies.slice
Restart=always
RestartSec=1
