# testing/btest/Baseline/tools.test-paths/..dir2.zeek-worker@.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Worker %i
SourcePath=<...>/config2
After=zeek-setup.service
After=zeek-manager.service
After=zeek-logger@.service
After=zeek-proxy@.service
Requires=zeek-setup.service
PartOf=zeek.target
StartLimitIntervalSec=0

[Service]
SyslogIdentifier=zeek-worker-%i
Type=exec
Nice=0
MemoryMax=
User=zeek
Group=zeek
WorkingDirectory=base-dir<...>/worker-%i
ReadWritePaths=base-dir<...>/worker-%i 
CapabilityBoundingSet=CAP_NET_RAW
AmbientCapabilities=CAP_NET_RAW
Environment=PATH=put-this:first:in-path:base-dir/bin:<...>/bin:<...>/bin:/bin
Environment=ZEEKPATH=put-this:first:in-zeekpath:base-dir<...>/generated-scripts:base-dir<...>/zeek:base-dir<...>/policy:base-dir<...>/site:base-dir<...>/builtin-plugins
Environment=CLUSTER_NODE=worker-%i
ExecStart=base-dir<...>/zeek -i ${INTERFACE} policy<...>/systemd-generator local frameworks<...>/connect
Slice=zeek-workers.slice
Restart=always
RestartSec=1
