# testing/btest/Baseline/tools.systemd-generator.test-zeek-systemd-generator-defaults/normal-dir.zeek-logger@.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Logger %i
SourcePath=<...>/config1
After=zeek-setup.service
Requires=zeek-setup.service
PartOf=zeek.target
StartLimitIntervalSec=0

[Service]
Type=exec
Nice=0
MemoryMax=
User=zeek
Group=zeek
WorkingDirectory=<...>/logger-%i
ReadWritePaths=<...>/logger-%i <...>/var 
Environment=PATH=<...>/bin:<...>/bin:<...>/bin:/bin
Environment=ZEEKPATH=<...>/generated-scripts:<...>/zeek:<...>/policy:<...>/site:<...>/builtin-plugins
Environment=CLUSTER_NODE=logger-%i
ExecStart=<...>/zeek policy<...>/systemd-generator local frameworks<...>/connect
Slice=zeek-loggers.slice
Restart=always
RestartSec=1
