# testing/btest/Baseline/tools.systemd-generator.test-zeek-systemd-generator-defaults/normal-dir.zeek-manager.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Manager
SourcePath=<...>/config1
After=zeek-setup.service
After=zeek-logger@.service
StopPropagatedFrom=zeek.target
StartLimitIntervalSec=0

[Service]
SyslogIdentifier=zeek-manager
Type=exec
Nice=0
MemoryMax=
User=zeek
Group=zeek
WorkingDirectory=<...>/manager
ReadWritePaths=<...>/manager 
Environment=PATH=<...>/bin:<...>/bin:<...>/bin:/bin
Environment=ZEEKPATH=<...>/generated-scripts:<...>/zeek:<...>/policy:<...>/site:<...>/builtin-plugins
Environment=CLUSTER_NODE=manager
ExecStart=<...>/zeek policy<...>/systemd-generator local frameworks<...>/zeromq
Slice=zeek-manager.slice
Restart=always
RestartSec=1
