# testing/btest/Baseline/tools.systemd-generator.test-zeek-systemd-generator-defaults/normal-dir.zeek-setup.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Setup
SourcePath=<...>/config1
PartOf=zeek.target

[Service]
Type=oneshot
Nice=0
MemoryMax=
User=
Group=
WorkingDirectory=
ExecStart=mkdir -p <...>/generated-scripts
ExecStart=<...>/zeek-cluster-layout-generator -L 1 -P 1 -W 1 -p 27760 -a 127.0.0.1 -m 9991 -b 0.0.0.0 -o <...>/cluster-layout.zeek
ExecStart=mkdir -p <...>/logs
ExecStart=chown zeek:zeek <...>/zeek
ExecStart=mkdir -p <...>/log-queue
ExecStart=chown zeek:zeek <...>/log-queue
ExecStart=mkdir -p <...>/manager
ExecStart=chown zeek:zeek <...>/manager
ExecStart=mkdir -p <...>/logger-1
ExecStart=chown zeek:zeek <...>/logger-1
ExecStart=mkdir -p <...>/proxy-1
ExecStart=chown zeek:zeek <...>/proxy-1
ExecStart=mkdir -p <...>/worker-1
ExecStart=chown zeek:zeek <...>/worker-1
RemainAfterExit=yes
