# testing/btest/Baseline/tools.systemd-generator.test-zeek-systemd-generator-defaults/normal-dir.zeek-proxy@.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Proxy %i
SourcePath=<...>/config1
After=zeek-setup.service
After=zeek-logger@.service
Requires=zeek-setup.service
PartOf=zeek.target

[Service]
Type=exec
Nice=0
MemoryMax=
User=zeek
Group=zeek
WorkingDirectory=<...>/proxy-%i
ReadWritePaths=<...>/proxy-%i 
Environment=PATH=<...>/bin:<...>/bin:<...>/bin:/bin
Environment=ZEEKPATH=<...>/generated-scripts:<...>/zeek:<...>/policy:<...>/site:<...>/builtin-plugins
Environment=CLUSTER_NODE=proxy-%i
ExecStart=<...>/zeek policy<...>/systemd-generator local frameworks<...>/connect
Slice=zeek-proxies.slice
Restart=always
RestartSec=1
