# testing/btest/Baseline/tools.systemd-generator.test-zeek-systemd-generator-defaults/normal-dir.zeek-archiver.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Archiver
SourcePath=<...>/config1
After=zeek-setup.service
PartOf=zeek.target
StartLimitIntervalSec=0

[Service]
Type=exec
Nice=0
MemoryMax=
User=zeek
Group=zeek
WorkingDirectory=
ReadWritePaths=<...>/log-queue <...>/zeek 
ExecStart=<...>/zeek-archiver  <...>/log-queue <...>/zeek
Slice=zeek-archiver.slice
Restart=always
RestartSec=1
