# testing/btest/Baseline/tools.systemd-generator.test-multi-interface/..dir1.zeek-worker-eth1@.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Worker %i (eth1)
SourcePath=<...>/config1
After=zeek-setup.service
After=zeek-manager.service
After=zeek-logger@.service
After=zeek-proxy@.service
StopPropagatedFrom=zeek.target
StartLimitIntervalSec=0

[Service]
SyslogIdentifier=zeek-worker-eth1-%i
Type=exec
User=zeek
Group=zeek
WorkingDirectory=<...>/worker-eth1-%i
ReadWritePaths=<...>/worker-eth1-%i 
Nice=-2
MemoryMax=1G
CapabilityBoundingSet=CAP_NET_RAW
AmbientCapabilities=CAP_NET_RAW
Environment=PATH=<...>/bin:<...>/bin:<...>/bin:/bin
Environment=ZEEKPATH=<...>/generated-scripts:<...>/zeek:<...>/policy:<...>/site:<...>/builtin-plugins
ExecStart=<...>/zeek -i ${INTERFACE} policy<...>/systemd-generator local AF_Packet::fanout_id=42 ignore_checksums=T frameworks<...>/zeromq
Slice=zeek-workers.slice
Restart=always
RestartSec=1
