# testing/btest/Baseline/tools.systemd-generator.test-multi-interface/..dir1.zeek-setup.service
### BTest baseline data generated by btest-diff. Do not edit. Use "btest -U/-u" to update. Requires BTest >= 0.63.
# Auto-generated, do not edit. Use drop-in files instead!
[Unit]
Description=Zeek Setup
SourcePath=<...>/config1
PartOf=zeek.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
User=
Group=
WorkingDirectory=
ExecStart=mkdir -p <...>/generated-scripts
ExecStart=<...>/zeek-cluster-layout-generator -L 3 -P 5 -W eth0:2,eth1:2,eth2:4 -p 27760 -a 127.0.0.1 -m 9991 -b 127.0.0.1 -o <...>/cluster-layout.zeek
ExecStart=mkdir -p <...>/zeek
ExecStart=chown zeek:zeek <...>/zeek
ExecStart=mkdir -p <...>/log-queue
ExecStart=chown zeek:zeek <...>/log-queue
ExecStart=mkdir -p <...>/manager
ExecStart=chown zeek:zeek <...>/manager
ExecStart=mkdir -p <...>/logger-1
ExecStart=chown zeek:zeek <...>/logger-1
ExecStart=mkdir -p <...>/logger-2
ExecStart=chown zeek:zeek <...>/logger-2
ExecStart=mkdir -p <...>/logger-3
ExecStart=chown zeek:zeek <...>/logger-3
ExecStart=mkdir -p <...>/proxy-1
ExecStart=chown zeek:zeek <...>/proxy-1
ExecStart=mkdir -p <...>/proxy-2
ExecStart=chown zeek:zeek <...>/proxy-2
ExecStart=mkdir -p <...>/proxy-3
ExecStart=chown zeek:zeek <...>/proxy-3
ExecStart=mkdir -p <...>/proxy-4
ExecStart=chown zeek:zeek <...>/proxy-4
ExecStart=mkdir -p <...>/proxy-5
ExecStart=chown zeek:zeek <...>/proxy-5
ExecStart=mkdir -p <...>/worker-eth0-1
ExecStart=chown zeek:zeek <...>/worker-eth0-1
ExecStart=mkdir -p <...>/worker-eth0-2
ExecStart=chown zeek:zeek <...>/worker-eth0-2
ExecStart=mkdir -p <...>/worker-eth1-1
ExecStart=chown zeek:zeek <...>/worker-eth1-1
ExecStart=mkdir -p <...>/worker-eth1-2
ExecStart=chown zeek:zeek <...>/worker-eth1-2
ExecStart=mkdir -p <...>/worker-eth2-1
ExecStart=chown zeek:zeek <...>/worker-eth2-1
ExecStart=mkdir -p <...>/worker-eth2-2
ExecStart=chown zeek:zeek <...>/worker-eth2-2
ExecStart=mkdir -p <...>/worker-eth2-3
ExecStart=chown zeek:zeek <...>/worker-eth2-3
ExecStart=mkdir -p <...>/worker-eth2-4
ExecStart=chown zeek:zeek <...>/worker-eth2-4
RemainAfterExit=yes
